# systemd/rqbit.path
[Unit]
Description=Automatically start rqbit after downloading a torrent file

[Path]
PathExistsGlob=$XDG_DOWNLOAD_DIR/*.torrent

[Install]
WantedBy=paths.target